# etc/systemd/system/mysql.service.d/override.conf
[Service]
ExecStartPre=/usr/share/mysql/mh_mysql-systemd-start pre
ExecStartPost=/usr/share/mysql/mh_mysql-systemd-start post
RestartSec=1
LimitNOFILE=1000000
LimitMEMLOCK=infinity
